# supervisor.ini
[inet_http_server]
port = :9001

[program:rig]
directory = /Users/bingoobjca/GitHub/rig
command = ./rig -u -o=false