# supervisor.ini
[inet_http_server]
port = :9001

[program:linx-test]
directory = /export/home/maguozhi/code/go/test_prometheus
command = ./test

[program:ddc]
directory = /export/home/maguozhi/unrealddc
command = ./unreal_httpddc_server
stdout_logfile = /export/home/maguozhi/unrealddc/logs/console.log
stderr_logfile = /export/home/maguozhi/unrealddc/logs/console.log